# run/com.austincaddell.ig-approval-worker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.austincaddell.ig-approval-worker</string>

    <key>ProgramArguments</key>
    <array>
      <string>/usr/local/bin/node</string>
      <string>/Users/austincaddell/.openclaw/workspace/ig-approval-automation/scripts/run-worker-loop.js</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/austincaddell/.openclaw/workspace/ig-approval-automation</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/austincaddell/.openclaw/workspace/ig-approval-automation/logs/worker.launchd.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/austincaddell/.openclaw/workspace/ig-approval-automation/logs/worker.launchd.err.log</string>
  </dict>
</plist>
